# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[inet_http_server]
port=127.0.0.1:9001

[supervisord]
logfile=/var/log/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[ctlplugin:quick]
supervisor.ctl_factory = supervisor_quick:make_quick_controllerplugin

[program:httpd]
command=/usr/sbin/httpd -c "ErrorLog /dev/stdout" -DFOREGROUND
redirect_stderr=true

[program:naemon]
command=/usr/bin/naemon /etc/naemon/naemon.cfg

[program:npcd]
command=/usr/sbin/npcd -f /etc/pnp4nagios/npcd.cfg

[group:naemon]
programs=httpd,naemon,npcd